# claude-daemon/com.naoyafurudono.claude-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.naoyafurudono.claude-daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/naoya-furudono/.local/bin/claude-daemon</string>
		<string>loop</string>
		<string>30m</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>StandardOutPath</key>
	<string>/Users/naoya-furudono/.claude/daemon/logs/launchd-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/naoya-furudono/.claude/daemon/logs/launchd-stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/naoya-furudono/.local/bin:/Users/naoya-furudono/.local/mise/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
